# learn-puppet/src/main/resources/modules/liferay-tomcat/templates/liferay_solaris_service.xml.erb
<?xml version="1.0"?>
<!DOCTYPE service_bundle SYSTEM "/usr/share/lib/xml/dtd/service_bundle.dtd.1">

<service_bundle type='manifest' name='Liferay Intra'>
<service name='<%= liferay_service_name %>' type='service' version='1'>

   <single_instance />

   <exec_method
      type='method'
      name='start'
      exec='<%= tomcat_home %>/bin/startup.sh'
      timeout_seconds='30' />

   <exec_method
      type='method'
      name='stop'
      exec='<%= tomcat_home %>/bin/shutdown.sh'
      timeout_seconds='30' />

   <instance name='default' enabled='false' />

   <stability value='Unstable' />

   <template>
      <common_name>
         <loctext xml:lang='C'>Liferay 6.1.1 GA2 Tomcat 7.0.27</loctext>
      </common_name>
      <documentation>
         <manpage title='tomcat' section='1' manpath='/usr/man' />
      </documentation>
   </template>

</service>
</service_bundle>
